# systemd/demo-portal-factory-reset.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/demo-portal-module/factory-reset-trigger
Unit=demo-portal-factory-reset.service

[Install]
WantedBy=paths.target
