# systemd/demo-portal-reboot.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/demo-portal-module/reboot-trigger
Unit=demo-portal-reboot.service

[Install]
WantedBy=paths.target
