# systemd/demo-portal-downgrade-start.service
[Unit]
Description=downgrade start sequence

[Service]
Type=oneshot
ExecStart=systemctl start aziot-identityd.service adu-agent.service
