# systemd/demo-portal-reboot.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/reboot/systemd-trigger
Unit=demo-portal-reboot.service

[Install]
WantedBy=paths.target
