# systemd/demo-portal-downgrade-stop.path

[Unit]
Before=paths.target

[Path]
PathChanged=/tmp/downgrade-stop
Unit=demo-portal-downgrade-stop.service

[Install]
WantedBy=paths.target
