# systemd/demo-portal-factory.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/factory-reset
Unit=demo-portal-factory.service

[Install]
WantedBy=paths.target
