# systemd/demo-portal-factory-reset.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/factory-reset/systemd-trigger
Unit=demo-portal-factory-reset.service

[Install]
WantedBy=paths.target
