# systemd/demo-portal-downgrade-start.path
[Unit]
Before=paths.target

[Path]
PathModified=/tmp/downgrade-start
Unit=demo-portal-downgrade-start.service

[Install]
WantedBy=paths.target
