# systemd/demo-portal-module.timer
[Unit]
Description=Restart demo-portal-module

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target