# systemd/demo-portal-factory.service
[Unit]
Description=factory reset

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/bin/factory_reset.sh /tmp/factory-reset'
