# systemd/icsdm-device-reboot.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/icsdm-device-service/reboot-trigger
Unit=icsdm-device-reboot.service

[Install]
WantedBy=paths.target
