# systemd/omnect-device-service.timer
[Unit]
Description=Restart omnect-device-service

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target
