# systemd/omnect-device-authorized_keys.path
[Unit]
Before=paths.target

[Path]
PathChanged=/run/omnect-authorized_keys

[Install]
WantedBy=paths.target
