# systemd/omnect-device-factory-reset.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/omnect-device-service/factory-reset-trigger
Unit=omnect-device-factory-reset.service

[Install]
WantedBy=paths.target
