# systemd/remote-features-factory-reset.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/remote-features-module/factory-reset-trigger
Unit=remote-features-factory-reset.service

[Install]
WantedBy=paths.target
