# systemd/remote-features-reboot.path

[Unit]
Before=paths.target

[Path]
PathChanged=/run/remote-features-module/reboot-trigger
Unit=remote-features-reboot.service

[Install]
WantedBy=paths.target
