# systemd/remote-features-module.timer
[Unit]
Description=Restart icsdm-device-service

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target
