# systemd/remote-features-module.timer
[Unit]
Description=Restart remote-features-module

[Timer]
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target
