# systemd/update-validation-observer.service
[Unit]
Description=Observe update validation timer service

[Service]
ExecStart=/usr/bin/update-validation-observer.sh