# systemd/update-validation-observer.timer
[Unit]
Description=Observe update validation timer

[Timer]
OnBootSec=10min

[Install]
WantedBy=timers.target
